# deploy/systemd/vw-egress-fw.service
# =============================================================================
# vw-egress-fw.service — persiste l'allow-list egress APRÈS Docker
# -----------------------------------------------------------------------------
# Pourquoi pas netfilter-persistent : il s'exécute AVANT Docker au boot ; la chaîne
# DOCKER-USER n'existe pas encore -> restauration échoue -> fenêtre sans filtrage.
# Cette unité, ordonnée After=docker.service, garantit l'application au bon moment
# et survit à un restart de Docker (script idempotent).
# -----------------------------------------------------------------------------
# Installation :
#   # Renseigner AUTHENTIK_IP dans vw-egress-fw.sh avant l'installation.
#   sudo cp deploy/firewall/vw-egress-fw.sh /usr/local/sbin/ && sudo chmod 700 /usr/local/sbin/vw-egress-fw.sh
#   sudo cp deploy/systemd/vw-egress-fw.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now vw-egress-fw.service
# =============================================================================
[Unit]
Description=Vaultwarden egress allow-list (DOCKER-USER)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vw-egress-fw.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
